# supervisord.conf
[supervisord]
nodaemon=true

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize --fpm-config /etc/php/7.0/fpm/php-fpm.conf

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n

[program:cron]
command=/usr/sbin/cron -f